# registration/conf/init.conf
# Upstart service configuration
APP=registration
env USER=$APP
env USER_HOME=/$APP-1.0-SNAPSHOT

env MANAGEMENT_PORT=9000
export MANAGEMENT_PORT

env LOGFILE=/$APP-1.0-SNAPSHOT/stdout.log

start on runlevel [2345]
stop on runlevel [016]

setuid $APP
chdir /$APP-1.0-SNAPSHOT

script
TOTAL_MEMORY=$(grep MemTotal /proc/meminfo | awk '{ print $2 }')
HEAP_SIZE_IN_MB=$(python -c "print int($TOTAL_MEMORY * 0.8 / 1024)")

JAVA_GC="-J-XX:NewRatio=1 \
-J-XX:+PrintGCDetails \
-J-XX:+PrintGCDateStamps \
-J-XX:+UseCompressedOops \
-J-XX:+UseConcMarkSweepGC \
-J-XX:+CMSClassUnloadingEnabled \
-J-verbose:gc \
-J-Xloggc:/$APP-1.0-SNAPSHOT/gc.log \
-J-XX:+UseGCLogFileRotation \
-J-XX:NumberOfGCLogFiles=5 \
-J-XX:GCLogFileSize=5M"

COMMAND="/$APP-1.0-SNAPSHOT/bin/$APP -J-Xmx${HEAP_SIZE_IN_MB}m -J-Xms${HEAP_SIZE_IN_MB}m $JAVA_GC"
echo "$COMMAND" >/$APP-1.0-SNAPSHOT/cmd.txt
$COMMAND >$LOGFILE 2>&1
end script